# oeb-studio-harness/ansible/roles/oeb_studio_harness_orchestrator/templates/oeb-studio-harness-backup.service.j2
[Unit]
Description=Studio Harness PostgreSQL backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/oeb-studio-harness-backup-postgres
User=root
StandardOutput=journal
StandardError=journal
